# root/etc/systemd/system/onBoot.service
[Unit]
Description=Custom setup on boot

[Service]
Type=oneshot
ExecStart=/opt/rootScripts/onBoot.sh

[Install]
WantedBy=multi-user.target